# taaza_usermanagement.service
[Unit]
Description=Taaza Usermanagement Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/henning
ExecStart=/home/henning/CTF/taaza_usermanagement
Restart=always
RestartSec=1s

[Install]
WantedBy=multi-user.target
