# taaza_usermanagement.service
[Unit]
Description=Taaza Usermanagement
After=network.target

[Service]
ExecStart=/home/CTF/taaza_usermanagement
WorkingDirectory=/home/CTF/taaza_usermanagement
Restart=always
RestartSec=3
User=root

[Install]
WantedBy=multi-user.target